# crontab
*/15 * * * * /usr/local/bin/python /app/speedtest_aio.py
